# system_files/etc/profile.d/xenia-env.sh
export OBS_VKCAPTURE=1
export PROTON_ENABLE_WAYLAND=1
